# unsorted/tablet-mode/squeekboard.service
[Unit]
BindsTo=tablet-mode.service

[Service]
ExecStart=/usr/bin/squeekboard

[Install]
WantedBy=tablet-mode.service
